# systemd/jtd.service
[Unit]
Description=jawn-tmux daemon
After=network.target

[Service]
Type=simple
ExecStart=/home/jamditis/.local/bin/jtd
Restart=always
RestartSec=5
Environment=HOME=/home/jamditis

[Install]
WantedBy=default.target
